# S05/nvda-tracker.service
[Unit]
Description=NVDA Stock Price Tracker
After=network.target

[Service]
Type=simple
User=nvda
Group=nvda
WorkingDirectory=/opt/nvda-tracker
ExecStart=/usr/local/bin/bun run /opt/nvda-tracker/nvda-tracker.ts
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
